# services-preview/fate-engine/scripts/setup/telegram-bot.service
[Unit]
Description=Fate-Engine Telegram Bot
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=lenovo
WorkingDirectory=/home/lenovo/.projects/fate-engine/services/telegram-service
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/home/lenovo/.projects/fate-engine-env/.env
ExecStart=/bin/bash -lc 'source /home/lenovo/.projects/fate-engine/.venv/bin/activate && python start.py bot'
Restart=always
RestartSec=5
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
